# contrib/com.annal.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.annal.server</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to match your install -->
        <string>/path/to/annal/.venv/bin/annal</string>
        <string>--transport</string>
        <string>streamable-http</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/annal.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/annal.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
